# systemd/input-leap.service
[Unit]
Description=Input Leap — Keyboard/Mouse Sharing
Documentation=https://github.com/input-leap/input-leap
Wants=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/mv/.local/input-leap/bin/input-leap
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target
