# systemd/input-leap.service
[Unit]
Description=Input Leap Server — Keyboard/Mouse Sharing
Documentation=https://github.com/input-leap/input-leap
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/home/mv/.local/input-leap/bin/input-leaps \
    --config /home/mv/.config/InputLeap/server.conf \
    --address :24800 \
    --no-daemon \
    --name mv-MS-7B16 \
    --use-x11 \
    --restart
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target
