# systemd/opencode-server.service
[Unit]
Description=OpenCode Pro Server (headless)
Documentation=https://opencode.ai/docs/server
After=network-online.target

[Service]
Type=simple
ExecStart=%h/.opencode/bin/opencode serve --port 4096 --hostname 127.0.0.1
Restart=on-failure
RestartSec=5
Environment=OPENCODE_ENABLE_EXA=1
Environment=OPENCODE_EXPERIMENTAL_LSP_TOOL=true
Environment=ANTHROPIC_API_KEY=dummy
Environment=ANTHROPIC_BASE_URL=http://127.0.0.1:3456
Environment=HOME=%h
WorkingDirectory=%h

# Security: localhost-only, no auth needed
# For remote access, add:
#   Environment=OPENCODE_SERVER_PASSWORD=your-strong-password
#   Environment=OPENCODE_SERVER_USERNAME=opencode

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
