# systemd/openvpn3-white-pc.service
# ⚠️ DEPRECATED — use openvpn3-session@white-pc (built-in)
# This file kept for reference only.
# Setup: openvpn3 config-import --config FILE.ovpn --name white-pc --persistent
#        openvpn3 config-acl --config white-pc --grant root --transfer-owner-session true
#        sudo systemctl enable --now openvpn3-session@white-pc
[Unit]
Description=OpenVPN3 White PC — DEPRECATED (use openvpn3-session@white-pc)
[Service]
ExecStart=/bin/true
[Install]
WantedBy=multi-user.target
